# mkosi/mkosi.extra/usr/lib/systemd/system/systemd-network.service.d/10-wants-nfsense-generator.conf
[Unit]
Wants=nfsense-generator.service
After=nfsense-generator.service
