# mkosi/mkosi.images/initrd/mkosi.extra/usr/lib/systemd/system/systemd-repart.service.d/10-override-seed.conf
# /opt/nfsense/mkosi/mkosi.extra/usr/lib/systemd/system/systemd-repart.service.d/10-override-seed.conf
[Service]
ExecStart=
ExecStart=/usr/bin/systemd-repart --seed=%m --dry-run=no
