# mkosi/mkosi.extra/usr/lib/systemd/system/unbound.service.d/10-runtime-config.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/unbound -d -c /usr/lib/unbound/unbound.conf

# Disable the unbound helper for chroot, resolveconf and rootkeys
# TODO Properly setup these in the base image
ExecStartPre=
ExecStopPost=
