# mkosi/mkosi.extra/usr/lib/systemd/system/nftables.service.d/10-runtime-config.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/nft -f /run/nfsense/nftables/nftables.conf
ExecReload=
ExecReload=/usr/sbin/nft -f /run/nfsense/nftables/nftables.conf
