# mkosi/mkosi.extra/usr/lib/systemd/system/chrony.service.d/10-runtime-config.conf
[Service]
ExecStart=
ExecStart=!/usr/sbin/chronyd -f /usr/lib/chrony/chrony.conf
