# mkosi/mkosi.extra/usr/lib/systemd/system/ulogd2.service.d/10-usr-config.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/ulogd --daemon --uid ulog --pidfile /run/ulog/ulogd.pid --configfile /usr/lib/ulogd/ulogd.conf
