# mkosi/mkosi.extra/usr/lib/systemd/system/kea-dhcp4-server.service.d/10-runtime-config.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/kea-dhcp4 -c /run/nfsense/kea/kea-dhcp4.conf
